# tools/agent-ssh-gateway/host/launchd/com.agentbot.queue-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.agentbot.queue-daemon.plist — jobs/incoming/ 自動執行 daemon

  部署：
    sudo cp host/bin/agent-queue-daemon.sh /usr/local/bin/
    sudo chmod 755 /usr/local/bin/agent-queue-daemon.sh
    cp host/launchd/com.agentbot.queue-daemon.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.agentbot.queue-daemon.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentbot.queue-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/agent-queue-daemon.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/ryan/meta-agent/tools/agent-ssh-gateway/logs/queue-daemon.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ryan/meta-agent/tools/agent-ssh-gateway/logs/queue-daemon.stderr.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ryan/meta-agent/tools/agent-ssh-gateway</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_SSH_GATEWAY_DIR</key>
    <string>/Users/ryan/meta-agent/tools/agent-ssh-gateway</string>
  </dict>
</dict>
</plist>
